# docker/supervisord.conf
; supervisord-Konfiguration für den Single-Container.
; Prozesse: mongod (optional), Node-Telefonie-Kern, Python-Admin-UI, Asterisk (optional).
; Optionale Prozesse werden über das entrypoint.sh per ENV ein-/ausgeschaltet
; (es schreibt die effektive Prozessliste; hier stehen die Defaults).

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:mongod]
; bind_ip 0.0.0.0: im Container auf allen Interfaces lauschen. Der Zugriff von außen wird
; bewusst nur über das Host-Loopback-Port-Mapping in run.sh freigegeben (-p 127.0.0.1:27017:27017).
command=/usr/bin/mongod --dbpath /data/db --bind_ip 0.0.0.0
autostart=%(ENV_SUPERVISOR_MONGOD)s
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:asterisk]
command=/usr/sbin/asterisk -f -C /etc/asterisk/asterisk.conf
autostart=%(ENV_SUPERVISOR_ASTERISK)s
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:voiceagent]
command=/usr/bin/node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:admin-ui]
command=/usr/bin/node /app/dist/admin/index.js
directory=/app
autostart=%(ENV_SUPERVISOR_ADMIN)s
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
